# piano-practice-tracker.service
[Unit]
Description=Piano Practice Tracker Service
After=sound.target network.target
Wants=sound.target

[Service]
Type=simple
User=aduston
WorkingDirectory=/home/aduston/pianolog
ExecStart=/home/aduston/pianolog/venv/bin/python /home/aduston/pianolog/main.py --prompt-each-session

# Automatic restart configuration
Restart=always
RestartSec=10

# Restart on any failure
StartLimitInterval=0

# Log output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=piano-tracker

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
